# remoteterm.service
[Unit]
Description=RemoteTerm for MeshCore
After=network.target

[Service]
Type=simple
User=remoteterm
Group=remoteterm
WorkingDirectory=/opt/remoteterm
ExecStart=/opt/remoteterm/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

# SPI mode reads data/config.yaml by default (MESHCORE_CONFIG_FILE). For USB LoRa,
# use scripts/manage_remoterm.sh to write /etc/remoterm/environment with MESHCORE_SERIAL_PORT.
Environment=MESHCORE_DATABASE_PATH=/opt/remoteterm/data/meshcore.db
EnvironmentFile=-/etc/remoterm/environment

# Serial USB: dialout. On Raspberry Pi OS also run: sudo usermod -aG spi,gpio remoteterm
SupplementaryGroups=dialout

[Install]
WantedBy=multi-user.target
